# ios/xpc_service_modern/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/@@PACKAGENAME@@.plist
